# examples/scheduling/systemd-timer/dynipupdate.service
[Unit]
Description=Dynamic IP DNS Updater
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=%i
EnvironmentFile=%h/.config/dynipupdate.env

# Pull latest image before running
ExecStartPre=/usr/bin/docker pull richleigh/dynipupdate:latest

# Run the update
ExecStart=/usr/bin/docker run --rm --network host --env-file %h/.config/dynipupdate.env richleigh/dynipupdate:latest

# Restart policy
Restart=on-failure
RestartSec=60s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dynipupdate

[Install]
WantedBy=default.target
