# examples/scheduling/com.user.dynipupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique label for this job -->
    <key>Label</key>
    <string>com.user.dynipupdate</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string>
        <string>run</string>
        <string>--rm</string>
        <string>--network</string>
        <string>host</string>
        <string>--env-file</string>
        <string>/Users/YOUR_USERNAME/.config/dynipupdate.env</string>
        <string>richleigh/dynipupdate:latest</string>
    </array>

    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Run on load (at login/boot) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/dynipupdate.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/dynipupdate.err.log</string>

    <!-- Keep the job alive even if it exits -->
    <key>KeepAlive</key>
    <false/>

    <!-- Run job even when battery power is low -->
    <key>LowPriorityIO</key>
    <true/>

    <!-- Environment variables (if not using --env-file) -->
    <!--
    <key>EnvironmentVariables</key>
    <dict>
        <key>CF_API_TOKEN</key>
        <string>your_token_here</string>
        <key>CF_ZONE_ID</key>
        <string>your_zone_id</string>
        <key>HOSTNAME</key>
        <string>your_hostname</string>
    </dict>
    -->
</dict>
</plist>
